# etc/profile.d/alias.sh
# Clear screen and scrollback in putty on clear or cls
alias clear="/bin/clear && printf '\033[3J'"
alias cls="clear"
