# etc/profile.d/stty-vim.sh
# Enable for vim Ctrl+Q shortcut
stty -ixon

